# module/remoteServer/nodeserviceOld.conf
#!upstart
description "node.js server"
author      "{url}"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn

    exec sudo -u root sh -c "/usr/bin/node {siteFolder}/server.js >> {log}/{url}.log 2>&1"
